# systemd/streamer.service
[Unit]
Description=YouTube Streamer Service
After=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/streamer.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
